# config/subscribe-hn-example.service
[Unit]
Description=SubscribeHN
Documentation=https://github.com/uint4096/subscribe-HN

[Service]
ExecStart=/usr/bin/subscribe_hn
User=<your_username>
Environment=TELOXIDE_TOKEN=<your_telegram_bot_id>
Environment=CHAT_ID=<your_telegram_chat_id>
Restart=always

[Install]
WantedBy=multi-user.target
